# samba/rootfs/etc/s6-overlay/s6-rc.d/wsdd/finish
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Keep the addon running when the optional wsdd service fails
# ==============================================================================

# shellcheck disable=SC2155
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"
readonly service="wsdd"

bashio::log.info \
  "Service ${service} exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

if [[ "${exit_code_service}" -eq 256 && "${exit_code_signal}" -eq 15 ]]; then
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo $((128 + exit_code_signal)) > /run/s6-linux-init-container-results/exitcode
  fi
  exec /run/s6/basedir/bin/halt
elif [[ "${exit_code_service}" -eq 256 ]]; then
  bashio::log.warning \
    "Optional service ${service} exited by signal ${exit_code_signal}; continuing without Windows discovery."
elif [[ "${exit_code_service}" -ne 0 ]]; then
  bashio::log.warning \
    "Optional service ${service} failed; continuing without Windows discovery."
fi
