# samba/rootfs/etc/s6-overlay/s6-rc.d/wsdd/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start wsdd service
# ==============================================================================
declare hostname
declare workgroup
declare -a wsdd_interfaces=()
declare -a interfaces=()

# shellcheck disable=SC1091  # Sourced from container filesystem path at runtime.
source /etc/s6-overlay/s6-rc.d/init-smbd/network-interfaces.sh

# Match Samba's advertised server name.
hostname=$(bashio::info.hostname)
if bashio::var.is_empty "${hostname}"; then
    bashio::log.warning "Can't read hostname, using default."
    hostname="hassio"
fi
workgroup=$(bashio::config 'workgroup')

get_enabled_interfaces interfaces

if [ ${#interfaces[@]} -eq 0 ]; then
    bashio::log.warning 'WSDD disabled: no supported interfaces found to bind on.'
    exec tail -f /dev/null
fi

# wsdd expects each interface to be preceded by "-i".
for interface in "${interfaces[@]}"; do
    wsdd_interfaces+=("-i" "${interface}")
done

bashio::log.info \
    "Starting the WSDD daemon on interfaces: $(printf '%s ' "${interfaces[@]}")" \
    "for ${workgroup}/${hostname}"

exec wsdd "${wsdd_interfaces[@]}" -n "${hostname}" -w "${workgroup}"
