# vlc/rootfs/etc/s6-overlay/s6-rc.d/vlc/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start VLC service
# ==============================================================================
# shellcheck disable=SC2155
readonly PASSWORD="$(cat /data/secret)"

# Send out discovery information to Home Assistant
/etc/s6-overlay/scripts/vlc-discovery &

# Run daemon
exec cvlc -I http --no-dbus --extraintf telnet --http-host 127.0.0.1 --http-password "INGRESS" --telnet-password "${PASSWORD}" --aout=pulse
