# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-post/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Post-start finalization of MariaDB service
# ==============================================================================

bashio::log.info "Deleting service information from Home Assistant"
bashio::services.delete "mysql" || true
