# samba/rootfs/etc/services.d/nmbd/run
#!/usr/bin/env bash
# ==============================================================================
# Start nmbd service
# ==============================================================================
exec nmbd \
  --foreground \
  --log-stdout \
  --no-process-group
