# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/bt-hci/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# bt_host_cpc_hci_bridge finish script
#==============================================================================
bashio::log.info "bt_host_cpc_hci_bridge ended with exit code ${1} (signal ${2})..."
