# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/bt-hci/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start bt_host_cpc_hci_bridge
# ==============================================================================

bashio::log.info "Starting bt_host_cpc_hci_bridge..."
cd /root
exec /usr/local/bin/bt_host_cpc_hci_bridge
