# homematic/rootfs/etc/services.d/lighttpd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start light httpd service
# ==============================================================================
bashio::log.info "Initialize webinterface routing..."
exec lighttpd-angel -D -f "/opt/hm/etc/lighttpd/lighttpd.conf"
