# nginx_proxy/rootfs/etc/periodic/daily/check_certificate_renewal
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Check certificate renewal
# ==============================================================================

CERTFILE=$(bashio::config 'certfile')

if [ ! -f /tmp/certificate_timestamp ] || [[ $(stat "/ssl/${CERTFILE}" -c %y) != $(cat /tmp/certificate_timestamp) ]]; then
  bashio::log.info "Reloading nginx to reload new certificate file"
  stat "/ssl/${CERTFILE}" -c %y > /tmp/certificate_timestamp
  nginx -c /etc/nginx.conf -s reload
fi
