# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/otbr-web/finish
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
#==============================================================================
# OpenThread BorderRouter web interface finish script
#==============================================================================
bashio::log.info "otbr-web ended with exit code ${1} (signal ${2})..."
