# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/otbr-web/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start OpenThread BorderRouter web interface
# ==============================================================================
bashio::log.info "Starting otbr-web..."

exec stdbuf -oL /usr/sbin/otbr-web -I wpan0 -d6 -s -p 8082 -a 127.0.0.1
