# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/otbr-web/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start OpenThread BorderRouter web interface
# ==============================================================================
# shellcheck source=./openthread_border_router/rootfs/etc/s6-overlay/scripts/otbr-agent-common
. /etc/s6-overlay/scripts/otbr-agent-common

bashio::log.info "Starting otbr-web..."
declare otbr_web_port
declare otbr_log_level_int

otbr_web_port="$(bashio::addon.port 8080)"
otbr_log_level_int="$(otbr_log_level_int)" \
    || bashio::exit.nok "Unknown otbr_log_level: $(bashio::config otbr_log_level)"

exec stdbuf -oL /usr/sbin/otbr-web -I wpan0 -d"${otbr_log_level_int}" -s -a "::" -p "${otbr_web_port}"
