# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/gdb-crash-dumper/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# GDB dumper finish script
#==============================================================================
bashio::log.info "gdb ended with exit code ${1} (signal ${2})..."
