# silabs-multiprotocol_dev/rootfs/etc/s6-overlay/s6-rc.d/universal-silabs-flasher/up
/etc/s6-overlay/scripts/universal-silabs-flasher-up
