# homematic/rootfs/etc/services.d/timesync/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start HomeMatic RF timesync service
# ==============================================================================
if bashio::config.false 'rf_enable'; then
    exec sleep infinity
fi

while true
do
    sleep 30m
    bashio::log.info "$(date '+%Y-%m-%d %H:%M:%S.%3N') Run SetInterfaceClock now."
    "${HM_HOME}/bin/SetInterfaceClock" 127.0.0.1:2001
done
