# silabs-multiprotocol_dev/rootfs/etc/s6-overlay/s6-rc.d/otbr-agent-rest-discovery/up
/etc/s6-overlay/scripts/otbr-agent-rest-discovery.sh
