# homematic/rootfs/etc/services.d/wired/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start HomeMatic wired service
# ==============================================================================
if bashio::config.false 'wired_enable'; then
    exec sleep infinity
fi

bashio::log.info "Starting wired..."
exec "${HM_HOME}/bin/hs485d" -g -i 0 -f /etc/config/hs485d.conf
