# almond/rootfs/etc/services.d/almond/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Almond service
# ==============================================================================
export THINGENGINE_BASE_URL

# Set the Ingress URL as Almond base URL for correct handling
THINGENGINE_BASE_URL=$(bashio::addon.ingress_entry)

mkdir -p /data/almond-server
if [ ! -f /data/almond-server/prefs.db ]; then
    # Skip authentication handling
    echo '{"server-login":{"password":"x","salt":"x","sqliteKeySalt":"x"}}' \
      > /data/almond-server/prefs.db
fi

# Send out discovery information to Home Assistant
./discovery &

# Register Home Assistant with Almond
./register &

# Start Almond
cd /opt/almond || bashio::exit.nok "Failed to change directory to Almond"
exec node main.js
