# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/mdns/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# mDNSResponder finish script
#==============================================================================
bashio::log.info "mDNS ended with exit code ${1} (signal ${2})..."
