# deconz/rootfs/etc/services.d/udev/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start udevd service
# ==============================================================================

if bashio::fs.directory_exists /run/udev; then
    exec sleep infinity
fi

# Triggers udev reload in the background
./trigger &

bashio::log.info "Starting udevd..."
exec /lib/systemd/systemd-udevd 2>/dev/null
