# vlc/rootfs/etc/services.d/vlc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start C-VLC daemon and enable discovery
# ==============================================================================
readonly PASSWORD="$(cat /data/secret)"

# Send out discovery information to Home Assistant
./discovery &

# Run daemon
exec cvlc -I http --extraintf telnet --http-host 127.0.0.1 --http-password "${PASSWORD}" --telnet-password "${PASSWORD}" --aout=pulse
