# letsencrypt/rootfs/etc/services.d/lets-encrypt/finish
#!/usr/bin/env bashio
# ==============================================================================
# Take down the S6 supervision tree when service is done 
# ==============================================================================
/run/s6/basedir/bin/halt

