# samba/rootfs/etc/s6-overlay/s6-rc.d/smbd/run
#!/command/with-contenv bashio
# ==============================================================================
# Start smbd service
# ==============================================================================
exec smbd \
  --foreground \
  --debug-stdout \
  --no-process-group
