# nginx_proxy/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bashio
# ==============================================================================
# Start crond service
# ==============================================================================

exec crond -f
