# nginx_proxy/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start crond service
# ==============================================================================

exec crond -f -d 9
