# piper/rootfs/etc/services.d/piper/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Piper service
# ==============================================================================
VOICE=$(bashio::config 'voice')
SPEAKER=$(bashio::config 'speaker')
LENGTH_SCALE=$(bashio::config 'length_scale')
NOISE_SCALE=$(bashio::config 'noise_scale')
NOISE_W=$(bashio::config 'noise_w')

exec python3 -m wyoming_piper --piper '/usr/share/piper/piper' --uri 'tcp://0.0.0.0:10200' --voice "${VOICE}" --speaker "${SPEAKER}" --length-scale "${LENGTH_SCALE}" --noise-scale "${NOISE_SCALE}" --noise-w "${NOISE_W}" --data-dir /data --download-dir /data
