# vlc/rootfs/etc/services.d/nginx/run
#!/usr/bin/env bashio
# ==============================================================================
# Start NGINX service
# ==============================================================================
exec nginx
