# silabs-multiprotocol_dev/rootfs/etc/s6-overlay/s6-rc.d/otbr-agent-configure/up
/etc/s6-overlay/scripts/otbr-agent-configure.sh
