# samba/rootfs/etc/s6-overlay/s6-rc.d/nmbd/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start nmbd service
# ==============================================================================

if ! bashio::config.true 'netbios'; then
    builtin exec sleep infinity
fi

builtin exec nmbd \
  --foreground \
  --debug-stdout \
  --no-process-group
