# samba/rootfs/etc/s6-overlay/s6-rc.d/nmbd/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start nmbd service
# ==============================================================================
exec nmbd \
  --foreground \
  --debug-stdout \
  --no-process-group
