# deconz/rootfs/etc/services.d/otau-deconz/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Download available firmware update for IKEA
# ==============================================================================

# Start OTA updates for deCONZ
bashio::log.info "Running the deCONZ OTA updater..."
deCONZ-otau-dl.sh &> /dev/null

exec sleep 259200
