# samba/rootfs/etc/services.d/smbd/run
#!/usr/bin/env bash
# ==============================================================================
# Start smbd service
# ==============================================================================
exec smbd \
  --foreground \
  --log-stdout \
  --no-process-group
