# speech_to_phrase/rootfs/etc/s6-overlay/s6-rc.d/speech-to-phrase/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start speech-to-phrase service
# ==============================================================================
cd /usr/src || exit 1

flags=()

if bashio::config.has_value 'hass_websocket_uri'; then
   hass_websocket_uri="$(bashio::config 'hass_websocket_uri')"
else
   hass_websocket_uri='ws://supervisor/core/websocket'
fi

if bashio::config.has_value 'hass_token'; then
   hass_token="$(bashio::config 'hass_token')"
else
   hass_token="${SUPERVISOR_TOKEN}"
fi

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

models_dir='/data/models'
mkdir -p "${models_dir}"

# shellcheck disable=SC2068
exec .venv/bin/python3 -m speech_to_phrase \
    --uri 'tcp://0.0.0.0:10300' \
    --tools-dir /usr/src/tools \
    --train-dir /share/speech-to-phrase/train \
    --models-dir "${models_dir}" \
    --retrain-on-start \
    --retrain-on-connect \
    --retrain-seconds 300 \
    --hass-token "${hass_token}" \
    --hass-websocket-uri "${hass_websocket_uri}" ${flags[@]}
