# zwave/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start mosquitto service for communication
# ==============================================================================
bashio::log.info "Starting internal MQTT broker..."

exec mosquitto -c /etc/mosquitto.conf
