# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/openocd/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# OpenOCD debugger finish script
#==============================================================================
bashio::log.info "openocd ended with exit code ${1} (signal ${2})..."
