# homematic/rootfs/etc/services.d/rfd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start HomeMatic RFD service
# ==============================================================================
if bashio::config.false 'rf_enable'; then
    exec sleep infinity
fi

bashio::log.info "Starting rfd..."
exec "${HM_HOME}/bin/rfd" -c -l 0 -f /etc/config/rfd.conf
