# samba/rootfs/etc/services.d/wsdd/run
#!/usr/bin/env bash
# ==============================================================================
# Start wsdd service
# ==============================================================================
exec wsdd2 -w
