# deconz/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start NGINX service
# ==============================================================================

# Wait for deCONZ to start before continuing
bashio::net.wait_for 40850

bashio::log.info "Starting Nginx..."
exec nginx
