# cec_scan/rootfs/etc/services.d/cec-scan/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start CEC scan service
# ==============================================================================

bashio::log.info "Starting CEC client scan..."
echo scan | cec-client -s -d 1
