# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-lock-core/finish
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# MariaDB client (to lock tables for backups) stopped
# ==============================================================================

bashio::log.info "MariaDB client (to lock tables for backups) stopped"
