# openwakeword/rootfs/etc/s6-overlay/s6-rc.d/openwakeword/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start openWakeWord service
# ==============================================================================
flags=()

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

# shellcheck disable=SC2068
exec python3 -m wyoming_openwakeword \
    --uri 'tcp://0.0.0.0:10400' \
    --zeroconf \
    --custom-model-dir /share/openwakeword \
    --threshold "$(bashio::config 'threshold')" \
    --trigger-level "$(bashio::config 'trigger_level')" ${flags[@]}
