# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/mdns/finish
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
#==============================================================================
# mDNSResponder finish script
#==============================================================================
bashio::log.info "mDNS ended with exit code ${1} (signal ${2})..."
