# matter_server/rootfs/etc/s6-overlay/s6-rc.d/matter-server/finish
/etc/s6-overlay/scripts/matter-server-finish
