# matter_server/rootfs/etc/s6-overlay/s6-rc.d/matter-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start matter-server service
# ==============================================================================
bashio::log.info "Starting Matter Server..."

declare server_port
declare log_level

log_level=$(bashio::string.lower "$(bashio::config log_level)")

# Bind to internal hassio network only unless user requests to expose
server_port="$(bashio::addon.port 5580)"
if ! bashio::var.has_value "${server_port}"; then
    server_port=5580
fi

# Reenable once Matter Server can bind to a specific IP
#if ! bashio::var.has_value "${server_port}"; then
#    export CHIP_WS_SERVER_HOST="$(bashio::addon.ip_address)"
#fi

# Send out discovery information to Home Assistant
/etc/s6-overlay/scripts/matter-server-discovery &

cd /root

exec matter-server --storage-path "/data" --port "${server_port}" \
                   --log-level "${log_level}" --vendorid 4939
