# matter_server/rootfs/etc/s6-overlay/s6-rc.d/matter-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start matter-server service
# ==============================================================================
bashio::log.info "Starting Matter Server..."

declare server_port
declare log_level
declare log_level_sdk
declare primary_interface
extra_args=()

if ! bashio::config.exists log_level; then
  bashio::log.magenta 'No log_level set in config, fallback to info'
fi
log_level=$(bashio::string.lower "$(bashio::config log_level info)")

# Make Matter SDK log level currently default to error
log_level_sdk=$(bashio::string.lower "$(bashio::config log_level_sdk error)")

if bashio::config.true "beta"; then
    bashio::log.info 'Upgrading Python Matter Server to latest pre-release'
    pip3 install --upgrade --pre python-matter-server[server]
fi

# Bind to internal hassio network only unless user requests to expose
server_port="$(bashio::addon.port 5580)"
if ! bashio::var.has_value "${server_port}"; then
    server_port=5580
    extra_args+=('--listen-address' "$(bashio::addon.ip_address)")
fi

if bashio::config.true "enable_test_net_dcl"; then
    extra_args+=('--enable-test-net-dcl')
fi

primary_interface="$(bashio::api.supervisor 'GET' '/network/info' '' 'first(.interfaces[] | select (.primary == true)) .interface')"

# Try fallback method (e.g. in case NetworkManager is not available)
if [ -z ${primary_interface} ]; then
  bashio::log.warning 'Trying fallback method to determine primary interface'
  primary_interface="$(ip --json route show default | jq --raw-output '.[0].dev')"
fi

if [ -z ${primary_interface} ] || [ ${primary_interface} == "null" ]; then
    bashio::exit.nok "No primary network interface found!"
fi

bashio::log.info "Using '${primary_interface}' as primary network interface."

# Send out discovery information to Home Assistant
/etc/s6-overlay/scripts/matter-server-discovery &

cd /root

if bashio::config.true "beta"; then
    exec /usr/bin/gdb --quiet -ex="set confirm off" -ex run -ex backtrace -ex "quit \$_exitcode" --args /usr/local/bin/python \
         /usr/local/bin/matter-server --storage-path "/data" --port "${server_port}" \
                       --log-level "${log_level}" --primary-interface "${primary_interface}" \
                       --paa-root-cert-dir "/data/credentials" \
                       --fabricid 2 --vendorid 4939 "${extra_args[@]}"
else
    exec /usr/local/bin/matter-server --storage-path "/data" --port "${server_port}" \
                       --log-level "${log_level}" --log-level-sdk "${log_level_sdk}" \
                       --primary-interface "${primary_interface}" \
                       --paa-root-cert-dir "/data/credentials" \
                       --fabricid 2 --vendorid 4939 "${extra_args[@]}"
fi
