# matter_server/rootfs/etc/s6-overlay/s6-rc.d/matter-server/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start matter-server service
# ==============================================================================
bashio::log.info "Starting Matter Server..."

declare server_port
declare log_level
declare primary_interface
declare install_target
declare matter_server_version
declare env_entry
declare env_name
declare env_value
declare env_vars_count
declare options
matter_server_args=()
extra_args=()

# Delete options removed in 9.0.0 so the Supervisor does not warn that they are
# no longer in the schema.
options=$(bashio::addon.options)
for old_key in log_level_sdk matter_sdk_wheels_version; do
  if bashio::jq.exists "${options}" ".${old_key}"; then
    bashio::log.info "Removing obsolete '${old_key}' option"
    bashio::addon.option "${old_key}"
  fi
done

if ! bashio::config.exists log_level; then
  bashio::log.magenta 'No log_level set in config, fallback to info'
fi
log_level=$(bashio::string.lower "$(bashio::config log_level info)")

# Map legacy Python Matter Server log levels to their matter.js equivalents.
case "${log_level}" in
  verbose) log_level="debug" ;;
  warning) log_level="warn" ;;
  critical) log_level="fatal" ;;
esac

# Resolve which matter-server build to run. The base image ships a bundled
# matter-server under /app. A pinned version wins over the "latest" beta switch.
# matter.js uses 0.x-2.x versions; a major >= 3 is almost certainly a leftover
# Python Matter Server version, so it is ignored to avoid installing the wrong package.
install_target=""
if bashio::config.has_value "matter_server_version"; then
  matter_server_version=$(bashio::config 'matter_server_version')
  if [[ "${matter_server_version}" =~ ^[^0-9]*([0-9]+) ]] && (( BASH_REMATCH[1] >= 3 )); then
    bashio::log.warning \
      "Ignoring matter_server_version '${matter_server_version}': a major version >= 3 looks like a Python Matter Server version. The matter.js Matter Server uses 0.x-2.x versions; please update or clear this option."
  else
    install_target="matter-server@${matter_server_version}"
  fi
fi
if ! bashio::var.has_value "${install_target}" && bashio::config.true "beta"; then
  install_target="matter-server@latest"
fi

if bashio::var.has_value "${install_target}"; then
  bashio::log.info "Installing ${install_target} from npm over the bundled server..."
  # matter.js native modules may need compiling; the base image purges the
  # toolchain, so add it back here. A failure here must not abort startup.
  if ! { apt-get update \
      && apt-get install -y --no-install-recommends \
        python3 make gcc g++ libbluetooth-dev libudev-dev \
      && cd /app \
      && { bashio::log.info "Fetching ${install_target} from npmjs.org (this can take a while)..."; true; } \
      && npm install --omit=dev --foreground-scripts "${install_target}"; }; then
    bashio::log.warning \
      "Failed to install ${install_target} from npm. Starting the bundled Matter Server instead."
  else
    npm cache clean --force || true
    find /app/node_modules -type d -name "cjs" -path "*/@matter/*" -exec rm -rf {} + 2>/dev/null || true
    find /app/node_modules -type d -name "cjs" -path "*/@project-chip/*" -exec rm -rf {} + 2>/dev/null || true
  fi
fi

# Export custom environment variables before starting the server
if bashio::config.has_value "matter_server_env_vars"; then
  env_vars_count=$(bashio::config 'matter_server_env_vars | length')
  for (( i=0; i < env_vars_count; i++ )); do
    env_entry="$(bashio::config "matter_server_env_vars[${i}]")"

    if [[ "${env_entry}" != *"="* ]]; then
      bashio::exit.nok "Invalid matter_server_env_vars entry: ${env_entry} (expected KEY=VALUE)"
    fi

    env_name="${env_entry%%=*}"
    env_value="${env_entry#*=}"

    if ! [[ "${env_name}" =~ ^[A-Za-z_][A-Za-z0-9_]*$ ]]; then
      bashio::exit.nok "Invalid environment variable name in matter_server_env_vars: ${env_name}"
    fi

    export "${env_name}=${env_value}"
    bashio::log.info "Exported Matter Server environment variable: ${env_name}"
  done
fi

# Remove leftover migration backup data from a former beta Matter Server.
if [ -d /data/.migrations ]; then
  bashio::log.info "Deleting former migration backup data ... this can take some time ... be patient"
  rm -rf /data/.migrations
fi

# Bind to internal hassio network only unless user requests to expose
server_port="$(bashio::addon.port 5580)"
if ! bashio::var.has_value "${server_port}"; then
    server_port=5580
    extra_args+=('--listen-address' "$(bashio::addon.ip_address)")
    extra_args+=('--listen-address' '127.0.0.1')  # so the base-image healthcheck (localhost) reaches us
fi

if bashio::config.true "enable_test_net_dcl"; then
    extra_args+=('--enable-test-net-dcl')
    extra_args+=('--ota-provider-dir' "/config/updates")
fi

# Resolve whether to enable time sync for nodes with the TimeSynchronization
# cluster. The 'time_sync' option has three modes:
#   - on:   always enable (warn if the host clock is not NTP synchronised)
#   - off:  always leave disabled
#   - auto: enable only when the host clock is NTP synchronised (default)
time_sync_mode="$(bashio::string.lower "$(bashio::config 'time_sync' 'auto')")"
ntp_synchronized="$(bashio::api.supervisor 'GET' '/supervisor/info' '' '.ntp_synchronized')"
time_sync_enabled=false

case "${time_sync_mode}" in
  on)
    time_sync_enabled=true
    if [ "${ntp_synchronized}" != "true" ]; then
        bashio::log.warning \
          "Time sync is enabled but host NTP is not synchronised — time pushed to Matter devices may be inaccurate"
    fi
    ;;
  off)
    bashio::log.info "Time sync is disabled ('time_sync' set to off)."
    ;;
  *)
    if [ "${ntp_synchronized}" = "true" ]; then
        time_sync_enabled=true
        bashio::log.info "Time sync 'auto': host NTP is synchronised, enabling time sync."
    else
        bashio::log.info "Time sync 'auto': host NTP is not synchronised, leaving time sync disabled."
    fi
    ;;
esac

if [ "${time_sync_enabled}" = "true" ]; then
    # The --enable-time-sync flag was added to the Matter Server in 1.2.0 (see
    # matter-js/matterjs-server#440). The server uses commander, which exits on
    # unknown options, so only pass the flag when the running server is new
    # enough; otherwise the option would break startup. A pre-release such as
    # 1.2.0-alpha is treated as 1.2.0 by comparing the MAJOR.MINOR.PATCH core.
    time_sync_min="1.2.0"
    running_version=""
    if [ -f /app/node_modules/matter-server/package.json ]; then
        running_version="$(jq --raw-output '.version // empty' /app/node_modules/matter-server/package.json)"
    fi
    running_core="${running_version%%-*}"
    if bashio::var.has_value "${running_core}" \
        && [ "$(printf '%s\n%s\n' "${time_sync_min}" "${running_core}" | sort -V | head -n1)" = "${time_sync_min}" ]; then
        extra_args+=('--enable-time-sync')
    else
        bashio::log.warning \
          "Matter Server ${running_version:-(unknown)} does not support time sync (requires >= ${time_sync_min}); '--enable-time-sync' not passed. Update the Matter Server to enable this feature."
    fi
fi

primary_interface="$(bashio::api.supervisor 'GET' '/network/info' '' 'first(.interfaces[] | select (.primary == true)) .interface')"

# Try fallback method (e.g. in case NetworkManager is not available)
# shellcheck disable=SC2086
if [ -z ${primary_interface} ]; then
  bashio::log.warning 'Trying fallback method to determine primary interface'
  primary_interface="$(ip --json route show default | jq --raw-output '.[0].dev')"
fi

# shellcheck disable=SC2086
if [ -z ${primary_interface} ] || [ ${primary_interface} == "null" ]; then
    bashio::exit.nok "No primary network interface found!"
fi

if bashio::config.true "ble_proxy"; then
    extra_args+=('--ble-proxy')
    if bashio::config.has_value "bluetooth_adapter_id"; then
        bashio::log.warning "Ignoring 'bluetooth_adapter_id' option: 'ble_proxy' is enabled (mutually exclusive)."
    fi
elif bashio::config.has_value "bluetooth_adapter_id"; then
  # shellcheck disable=SC2207
  extra_args+=('--bluetooth-adapter' $(bashio::config 'bluetooth_adapter_id'))
fi

if bashio::config.has_value "matter_server_args"; then
  # shellcheck disable=SC2207
  extra_args+=($(bashio::config 'matter_server_args'))
fi

bashio::log.info "Using '${primary_interface}' as primary network interface."

# Send out discovery information to Home Assistant
/etc/s6-overlay/scripts/matter-server-discovery &

# shellcheck disable=SC2164
cd /root

# shellcheck disable=SC2206
matter_server_args+=(
  '--storage-path' "/data"
  '--port' "${server_port}"
  '--log-level' "${log_level}"
  '--primary-interface' "${primary_interface}"
  '--fabricid' 2
  '--vendorid' 4939
  ${extra_args[@]}
)

exec node --enable-source-maps /app/node_modules/matter-server/dist/esm/MatterServer.js "${matter_server_args[@]}"
