# dnsmasq/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/dnsmasq.conf"

# Run dnsmasq
bashio::log.info "Starting dnsmasq..."
exec dnsmasq -C "${CONFIG}" -z < /dev/null

