# dnsmasq/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/dnsmasq.conf"

# Run dnsmasq
bashio::log.info "Starting dnsmasq..."

# Set max open file limit to speed up startup
ulimit -n 1024

exec dnsmasq -C "${CONFIG}" -z < /dev/null

