# mosquitto/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs the NGINX daemon
# ==============================================================================
bashio::log.info "Starting NGINX for authentication handling..."
exec nginx
