# assist_satellite/rootfs/etc/s6-overlay/s6-rc.d/assist_satellite/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start the Assist Satellite service (linux-voice-assistant)
# ==============================================================================
declare -a extra_args=()

# Debug logging
if bashio::config.true 'debug_logging'; then
    extra_args+=('--debug')
fi

# Thinking sound (processing sound played while assistant thinks)
if bashio::config.true 'thinking_sound'; then
    extra_args+=('--enable-thinking-sound')
fi

# Optional friendly name (auto-generated from MAC if empty)
declare name
name="$(bashio::config 'name')"
if [[ -n "${name}" ]]; then
    extra_args+=('--name' "${name}")
fi

declare auto_gain_control
auto_gain_control = "$(bashio::config 'mic_auto_gain')"
if auto_gain_control; then
    extra_args+=('--mic-auto-gain',"${auto_gain_control}")
fi

declare noise_suppression
noise_suppression = "$(bashio::config 'mic_noise_suppression')"
if noise_suppression; then
    extra_args+=('--mic-noise-suppression',"${noise_suppression}")
fi

declare aec_enabled=$(bashio::config 'aec_enabled')
if [[ "${aec_enabled}" == "true" ]]; then
    pw-cli load-module libpipewire-module-echo-cancel \
        'aec.method=webrtc source.props={ node.name=aec_mic }' || bashio::log.error "Failed to load AEC"
    extra_args+=("--audio-input-device" "aec_mic")
fi

exec python3 -m linux_voice_assistant \
    --port "$(bashio::config 'port')" \
    --wake-model "$(bashio::config 'wake_word')" \
    --refractory-seconds "$(bashio::config 'refractory_seconds')" \
    --wake-word-dir /usr/src/wakewords \
    --wakeup-sound /usr/src/sounds/wake_word_triggered.flac \
    --timer-finished-sound /usr/src/sounds/timer_finished.flac \
    --processing-sound /usr/src/sounds/processing.wav \
    --mute-sound /usr/src/sounds/mute_switch_on.flac \
    --unmute-sound /usr/src/sounds/mute_switch_off.flac \
    --preferences-file /share/assist_satellite/preferences.json \
    --download-dir /share/assist_satellite/local \
    "${extra_args[@]}"