# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/otbr-web/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# OpenThread BorderRouter web interface finish script
#==============================================================================
bashio::log.info "otbr-web ended with exit code ${1} (signal ${2})..."
