# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/otbr-web/finish
/etc/s6-overlay/scripts/otbr-web-finish
