# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/otbr-web/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start OpenThread BorderRouter web interface
# ==============================================================================
bashio::log.info "Starting otbr-web..."
exec /usr/sbin/otbr-web -I wpan0 -d6 -p 8080
