# gammu/rootfs/etc/s6-overlay/s6-rc.d/gammu/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

export GSM_DEVICE="$(bashio::config 'gsm_device')"

# shellcheck disable=SC2068
exec python3 -m pysmsgateway
