# whisper/rootfs/etc/s6-overlay/s6-rc.d/whisper/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Whisper service
# ==============================================================================
flags=()

if [ "$(uname -m)" == "x86_64" ] && ! grep -qw 'avx' /proc/cpuinfo; then
    bashio::log.warning "Your CPU does not support Advanced Vector Extensions (AVX). Whisper will run slower than normal."
fi

model="$(bashio::config 'model')"
stt_library="$(bashio::config 'stt_library')"

if [ "${model}" = 'custom' ]; then
    # Override with custom model
    model="$(bashio::config 'custom_model')"
    if [ -z "${model}" ]; then
        bashio::exit.nok "Custom model is not set"
    fi
    if [ "${stt_library}" = 'auto' ]; then
        # Use existing option
        stt_library="$(bashio::config 'custom_model_type')"
    fi
elif [ "${model}" != 'auto' ]; then
    if [ "${stt_library}" = 'auto' ]; then
        # Default to faster whisper if model is selected
        stt_library='faster-whisper'
    fi
fi

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

# shellcheck disable=SC2068
exec python3 -m wyoming_faster_whisper \
    --uri 'tcp://0.0.0.0:10300' \
    --zeroconf \
    --model "${model}" \
    --stt-library "${stt_library}" \
    --beam-size "$(bashio::config 'beam_size')" \
    --language "$(bashio::config 'language')" \
    --initial-prompt "$(bashio::config 'initial_prompt')" \
    --data-dir /data \
    --download-dir /data ${flags[@]}
