# whisper/rootfs/etc/s6-overlay/s6-rc.d/whisper/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Whisper service
# ==============================================================================
exec python3 -m wyoming_faster_whisper \
    --uri 'tcp://0.0.0.0:10300' \
    --model "$(bashio::config 'model')" \
    --beam-size "$(bashio::config 'beam_size')" \
    --language "$(bashio::config 'language')" \
    --data-dir /data \
    --download-dir /data
