# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/cpcd/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Co-Processor Communication Daemon (CPCd)
# ==============================================================================

bashio::log.info "Starting cpcd..."
exec s6-notifyoncheck -d -s 300 -w 300 -n 0 \
    "/usr/bin/stdbuf" -o0 /usr/local/bin/cpcd
