# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/zigbeed/finish
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
#==============================================================================
# EmberZNet Zigbee Daemon finish script
#==============================================================================
bashio::log.info "zigbeed ended with exit code ${1} (signal ${2})..."
